# configs/jp/airootfs/etc/profile.d/fcitx5.sh
export GLFW_IM_MODULE=ibus
export XMODIFIERS=@im=fcitx
