# configs/termy/airootfs/etc/profile.d/termy.sh
export EDITOR=nvim
export SHELL=/usr/bin/nu
